# utils/vector/files/etc/uci-defaults/50_vector
#!/bin/sh

# Disable auto-start vector by default
ubus call service delete vector

exit 0